# packaging/launchd/com.ai-agent-owner.op-sa-broker.plist.tmpl
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>{{BROKER_LABEL}}</string>

  <key>ProgramArguments</key>
  <array>
    <string>{{BROKER_BINARY}}</string>
  </array>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <true/>

  <key>LimitLoadToSessionType</key>
  <array>
    <string>Aqua</string>
    <string>Background</string>
  </array>

  <key>WorkingDirectory</key>
  <string>{{WORKING_DIRECTORY}}</string>

  <key>StandardOutPath</key>
  <string>{{STDOUT_PATH}}</string>

  <key>StandardErrorPath</key>
  <string>{{STDERR_PATH}}</string>
</dict>
</plist>
